# scripts/astronomical-data.service
[Unit]
Description=Initialize Astronomical Data for MidnightRider
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=aneto
ExecStart=/home/aneto/docker/signalk/scripts/init-astronomical-data.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
